# packaging/deb/DEBIAN/postinst
#!/bin/bash
set -e

# postinst script for sysmedic
# see: dh_installdeb(1)

case "$1" in
    configure)
        # Create necessary directories
        mkdir -p /etc/sysmedic
        mkdir -p /var/lib/sysmedic

        # Set proper permissions
        chmod 755 /etc/sysmedic
        chmod 755 /var/lib/sysmedic

        # Create default config if it doesn't exist
        if [ ! -f /etc/sysmedic/config.yaml ]; then
            cat > /etc/sysmedic/config.yaml << 'EOF'
monitoring:
  check_interval: 60
  cpu_threshold: 80
  memory_threshold: 80
  persistent_time: 60

users:
  cpu_threshold: 80
  memory_threshold: 80
  persistent_time: 60

reporting:
  period: "hourly"
  retain_days: 30

email:
  enabled: false
  smtp_host: "smtp.gmail.com"
  smtp_port: 587
  tls: true

websocket:
  enabled: true
  port: 8060
  secret: ""

user_filtering:
  min_uid_for_real_users: 1000
  ignore_system_users: true
  min_cpu_percent: 5.0
  min_memory_percent: 5.0
  min_process_count: 1
  excluded_users: ["root", "daemon", "bin", "sys", "sync", "games", "man", "lp", "mail", "news", "uucp", "proxy", "www-data", "backup", "list", "irc", "gnats", "nobody", "systemd+", "syslog", "_apt"]
  included_users: []

data_path: "/var/lib/sysmedic"

user_thresholds: {}
EOF
            chmod 644 /etc/sysmedic/config.yaml
        fi

        # Reload systemd and enable services
        systemctl daemon-reload
        systemctl enable sysmedic.doctor.service
        systemctl enable sysmedic.websocket.service

        echo "SysMedic has been installed successfully!"
        echo "Configuration file: /etc/sysmedic/config.yaml"
        echo "Enable services with:"
        echo "  sudo systemctl enable sysmedic.doctor"
        echo "  sudo systemctl enable sysmedic.websocket"
        echo "Start services with:"
        echo "  sudo systemctl start sysmedic.doctor"
        echo "  sudo systemctl start sysmedic.websocket"
        echo "Or use the CLI:"
        echo "  sysmedic daemon start"
        echo "  sysmedic websocket start"
        echo "View logs with:"
        echo "  sudo journalctl -u sysmedic.doctor -f"
        echo "  sudo journalctl -u sysmedic.websocket -f"
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

exit 0
