# packaging/deb/DEBIAN/prerm
#!/bin/bash
set -e

# prerm script for sysmedic
# see: dh_installdeb(1)

case "$1" in
    remove|upgrade|deconfigure)
        # Stop the services if they're running
        if systemctl is-active --quiet sysmedic.doctor.service; then
            echo "Stopping sysmedic doctor service..."
            systemctl stop sysmedic.doctor.service
        fi

        if systemctl is-active --quiet sysmedic.websocket.service; then
            echo "Stopping sysmedic websocket service..."
            systemctl stop sysmedic.websocket.service
        fi

        # Disable the services if they're enabled
        if systemctl is-enabled --quiet sysmedic.doctor.service; then
            echo "Disabling sysmedic doctor service..."
            systemctl disable sysmedic.doctor.service
        fi

        if systemctl is-enabled --quiet sysmedic.websocket.service; then
            echo "Disabling sysmedic websocket service..."
            systemctl disable sysmedic.websocket.service
        fi
        ;;

    failed-upgrade)
        ;;

    *)
        echo "prerm called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

exit 0
